# trade_bot.service
[Unit]
Description=Trade Bot AI Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade_bot_ai
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/local/bin/python3.12 /root/trade_bot_ai/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

